# autowall.timer.d/10-timer.conf
[Timer]
OnBootSec=15m
OnUnitActiveSec=15m
