# autowall.timer.d/2-timer.conf
[Timer]
OnBootSec=30m
OnUnitActiveSec=30m
