# dynamicwall.timer
[Unit]
Description=dynamic-wall timer

[Timer]
OnBootSec=30m
OnUnitActiveSec=30m

[Install]
WantedBy=multi-user.target timers.target
