# config/systemd/user/dynamicwall.service
[Unit]
Description=dynamic-wall

[Service]
Type=oneshot
ExecStart=dynamicwall
