# config/systemd/user/dynamicwall.service
[Unit]
Description=dynamic-wall unit

[Service]
Type=simple
ExecStart=dynamicwall
RestartSec=3
Restart=on-failure

[Install]
WantedBy=default.target
