# config/systemd/user/dynamicwall.service
[Unit]
Description=dynamic-wall

[Service]
Type=oneshot
User=USERNAME
Group=USERNAME
ExecStart=/bin/bash /home/USERNAME/bin/dynamic-wall/dynamic
