# config/systemd/user/dynamicwall.service
[Unit]
Description=dynamic-wall
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=dynamicwall

[Install]
WantedBy=multi-user.target
