# config/systemd/user/dynamicwall.service
[Unit]
Description=dynamic-wall unit
After=graphical.target

[Service]
Type=oneshot
ExecStart=dynamicwall

[Install]
WantedBy=multi-user.target
