# config/systemd/user/dynamicwall.timer
[Unit]
Description=dynamic-wall timer

[Timer]
OnBootSec=45m
OnUnitActiveSec=45m

[Install]
WantedBy=multi-user.target timers.target
