# config/systemd/user/dynamicwall.service
[Unit]
Description=dynamic-wall

[Service]
Type=oneshot
ExecStart=/bin/bash /home/USERNAME/bin/dynamic-wall/dynamic
