# config/systemd/user/dynamicwall.timer.d/2-timer.conf
[Timer]
OnBootSec=
OnUnitActiveSec=
OnBootSec=30m
OnUnitActiveSec=30m
